# core-xprojects/webrtc/build/src/third_party/pipewire/src/daemon/systemd/user/pipewire.socket
[Unit]
Description=Multimedia System

[Socket]
Priority=6
ListenStream=%t/pipewire-0

[Install]
WantedBy=sockets.target
