# core-xprojects/webrtc/build/src/third_party/pipewire/src/daemon/systemd/user/pipewire-pulse.socket
[Unit]
Description=PipeWire PulseAudio
ConditionUser=!root

[Socket]
Priority=6
ListenStream=%t/pulse/native

[Install]
WantedBy=sockets.target
