# core-xprojects/webrtc/build/src/third_party/pipewire/src/daemon/systemd/system/pipewire.socket
[Unit]
Description=Multimedia System

[Socket]
Priority=6
ListenStream=%t/pipewire/pipewire-0
SocketUser=pipewire
SocketGroup=pipewire
SocketMode=0660

[Install]
WantedBy=sockets.target
